# supervisor.conf
[program:queue]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --tries=10 --queue=default --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue.log

[program:queue_cashier]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=cashier --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_cashier.log

[program:monitoring_closing]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=monitoring_closing --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_monitoring_closing.log

[program:monitoring_closing_summary]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=monitoring_closing_summary --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_monitoring_closing_summary.log

[program:report_transaction]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=report_transaction --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_report_transaction.log

[program:forecast_conversion]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=forecast_conversion --tries=3  --timeout=240 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_forecast_conversion.log

[program:forecast_conversion_2]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=forecast_conversion_2 --tries=3  --timeout=240 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_forecast_conversion_2.log

[program:forecast_conversion_detail]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=forecast_conversion_detail --tries=5  --timeout=240 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_forecast_conversion_detail.log

[program:notification]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=notification --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/queue_notification.log

[program:po_supplier]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=po_supplier --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/po_supplier.log

[program:forecast_import]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=forecast_import --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/forecast_import.log

[program:so_import]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=so_import --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/so_import.log

[program:so_stock]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=so_stock --tries=3  --timeout=0 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/so_stock.log

[program:report_recipe]
command=php /var/www/html/api.erpkhasanahsari.com/artisan queue:work --queue=report_recipe --tries=3  --timeout=3600 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/api.erpkhasanahsari.com/storage/logs/report_recipe.log
